# iptv-api/docker/crontab
# ============================================
# Crontab para sincronización automática
# Ejecutar sync_iptv.py cada 2 horas
# ============================================

# Cada 2 horas en punto (00:00, 02:00, 04:00, etc.)
0 */2 * * * python /app/scripts/sync_iptv.py >> /var/log/sync.log 2>&1

# Alternative: Cada 2 horas desde el inicio del contenedor
# */2 * * * python /app/scripts/sync_iptv.py >> /var/log/sync.log 2>&1
